# supervisord.conf
[supervisorctl]
[program:twisted]
command = bin/twistd -ny twemail.tac
directory = /home/dotcloud/
