# supervisord.conf
[program:twisted]
command = twistd -ny twemail.tac
