# supervisord.conf
[program:twisted]
command = twistd -ny twemail.tac
directory = /home/dotcloud/current/
